# deploy/project-flomo-server.service
[Unit]
Description=Project Flomo Node Server
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/project-flomo
Environment=NODE_ENV=production
Environment=HOST=127.0.0.1
Environment=PORT=3000
EnvironmentFile=-/etc/project-flomo.env
ExecStart=/usr/bin/env node server/index.mjs
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
